# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:php-fpm]
; Asegúrate de que el comando sea correcto para la imagen base
command=/usr/local/sbin/php-fpm --fpm-config /usr/local/etc/php-fpm.conf -F
autostart=true
autorestart=true
